# launchd/org.datamak.check-main-db-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.datamak.check-main-db-backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/.venv/bin/python</string>
    <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/tools/backup_db_to_flux.py</string>
    <string>check</string>
    <string>--status-file</string>
    <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/output/flux_db_backup/backup_status.json</string>
    <string>--max-age-hours</string>
    <string>192</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>DTWIN_ROOT</key>
    <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak</string>

  <key>StartInterval</key>
  <integer>21600</integer>

  <key>StandardOutPath</key>
  <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/output/flux_db_backup/check.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/output/flux_db_backup/check.stderr.log</string>
</dict>
</plist>
