# launchd/org.datamak.backup-main-db-to-flux.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.datamak.backup-main-db-to-flux</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/.venv/bin/python</string>
    <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/tools/backup_db_to_flux.py</string>
    <string>backup</string>
    <string>--db</string>
    <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/gyrokinetic_simulations.db</string>
    <string>--status-file</string>
    <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/output/flux_db_backup/backup_status.json</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>DTWIN_ROOT</key>
    <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/output/flux_db_backup/backup.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jdominsk/Documents/Projects/AIML_database/Datamak/output/flux_db_backup/backup.stderr.log</string>
</dict>
</plist>
